# browser/files/rc.local
#!/bin/sh
rm -f /home/pi/.config/chromium/Singleton*
exit 0
